# files/default/nix-worker.conf
start on runlevel [2345]
stop  on runlevel [016]
respawn
exec /nix/var/nix/profiles/default/bin/nix-worker --daemon
